# photoframe/etc/systemd/system/photos-on.service
[Unit]
Description=Power on photos
After=network.target

[Service]
User=pi
ExecStart=/usr/local/bin/start-photos.sh
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
WorkingDirectory=/home/pi
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
